# riak/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid

[program:riak]
command=bash -c "/usr/sbin/riak console"
numprocs=1
startsecs=10
autostart=true
autorestart=false
priority=1
user=riak
redirect_stderr=true
environment=HOME="/var/lib/riak"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[eventlistener:riak_exit]
command=/usr/bin/kill_supervisor.py
process_name=riak
events=PROCESS_STATE_FATAL,PROCESS_STATE_EXITED
